# dkms-gentoo/dkms
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount net
	keyword -timeout
}

start() {
	ebegin "Starting dkms"
	  dkms-gentoo --check
	eend $?
}

stop() {
	ebegin "Stopping dkms"
	  dkms-gentoo --db
	eend $?
}


